# rclone-systemd/mnt-idrive-shows.automount
[Unit]
Description=AutoMount for /mnt/idrive/shows

[Automount]
Where=/mnt/idrive/shows
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
