# rclone-systemd/mnt-idrive-requests.automount
[Unit]
Description=AutoMount for /mnt/idrive/requests

[Automount]
Where=/mnt/idrive/requests
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
