# rclone-systemd/mnt-idrive-games.automount
[Unit]
Description=AutoMount for /mnt/idrive/games

[Automount]
Where=/mnt/idrive/games
TimeoutIdleSec=600

[Install]
WantedBy=multi-user.target
